# home/dot_config/autostart/blueman-applet.desktop
[Desktop Entry]
Type=Application
Name=Blueman Applet
Exec=blueman-applet
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true