# home/dot_config/autostart/waybar.desktop
[Desktop Entry]
Type=Application
Name=Waybar
Exec=waybar
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true